# dot_config/systemd/user/amixer-autounmute.service
[Unit]
Description=Disable Auto-Mute on sound card 3

[Service]
Type=oneshot
ExecStart=/usr/bin/amixer -c 3 sset "Auto-Mute Mode" Disabled
