# dot_config/systemd/user/kioworker-killer.service
[Unit]
Description=Kill any long living kioworker threads

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill kioworker
